# filebrowser/filebrowser.ini
[program:filebrowser]
command=/usr/local/bin/filebrowser -r /var/www/html/storage --address 0.0.0.0 --port 8080
autostart=true
autorestart=true
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser_err.log
